# app/app-gpu/gpu-tools/deploy/systemd/opus-gpu-watchdog.service
[Unit]
Description=OPUS-GPU Watchdog Service
Documentation=https://github.com/opus-gpu/app-gpu
After=opus-gpu.service
Requires=opus-gpu.service

[Service]
Type=simple
User=miner
Group=miner

# Working directory
WorkingDirectory=/opt/opus-gpu

# Watchdog executable
ExecStart=/usr/local/bin/gpu-watchdog --interval 10s --prometheus-url http://localhost:9090

# Restart policy
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=opus-gpu-watchdog

# Resource limits
MemoryMax=256M
CPUQuota=50%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
