# app/app-gpu/gpu-tools/deploy/systemd/opus-gpu.service
[Unit]
Description=OPUS-GPU High-Performance Cryptocurrency Miner
Documentation=https://github.com/opus-gpu/app-gpu
After=network-online.target
Wants=network-online.target
# If using NVIDIA GPUs, ensure drivers are loaded
After=nvidia-persistenced.service
Requires=nvidia-persistenced.service

[Service]
Type=simple
User=miner
Group=miner

# Working directory
WorkingDirectory=/opt/opus-gpu

# Main executable
ExecStart=/usr/local/bin/gpu-miner --config /etc/opus-gpu/app.toml

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=5min
StartLimitBurst=5

# Environment variables
Environment="RUST_LOG=info"
Environment="CONFIG_PATH=/etc/opus-gpu/app.toml"
Environment="LOG_PATH=/var/log/opus-gpu"
Environment="CUDA_VISIBLE_DEVICES=0,1"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=opus-gpu

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=4G
CPUQuota=400%  # 4 cores max

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/opus-gpu /opt/opus-gpu/data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
PrivateMounts=true
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
SystemCallErrorNumber=EPERM

# Capabilities (minimal set for GPU access)
AmbientCapabilities=
CapabilityBoundingSet=

# GPU device access
DeviceAllow=/dev/nvidia0 rw
DeviceAllow=/dev/nvidia1 rw
DeviceAllow=/dev/nvidiactl rw
DeviceAllow=/dev/nvidia-uvm rw
DeviceAllow=/dev/nvidia-uvm-tools rw

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
